# files/systemd/clamav/freshclam.service
[Unit]
Description=updates clamav virus databases
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=freshclam
