# files/systemd/tealdeer/tealdeer.service
[Unit]
Description=updates tldr local cache
Requires=network-online.target
After=network-online.target

[Service]
ExecStart=tldr --update
