# files/systemd/tealdeer/tealdeer.timer
[Unit]
Description=updates tldr local cache

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=3hours

[Install]
WantedBy=multi-user.target
